# common/src/stack/aws-server/aws-server-init.service
[Unit]
After=network.target

[Service]
ExecStart=/opt/stack/sbin/aws-server-init

[Install]
WantedBy=default.target