# common/src/stack/aws-server/aws-server-init.service
[Unit]
After=network.target
After=mariadb.service 

[Service]
ExecStart=/opt/stack/sbin/aws-server-init

[Install]
WantedBy=multi-user.target
