# common/src/stack/config/etc/systemd/system/last.target
[Unit]
Description=Run Services Last
After=multi-user.target
Conflicts=shutdown.target

[Install]
WantedBy=multi-user.target
