# common/src/stack/docker/docker-barnacle.service
[Unit]
After=last.target

[Service]
ExecStart=/opt/stack/sbin/docker-barnacle

[Install]
WantedBy=default.target