# common/src/stack/aws-client/aws-client-register.service
[Unit]
After=network.target

[Service]
ExecStart=/opt/stack/sbin/aws-client-register

[Install]
WantedBy=default.target