# common/src/stack/barnacle/aws-barnacle.service
[Unit]
After=network.target

[Service]
ExecStart=/opt/stack/sbin/aws-barnacle

[Install]
WantedBy=multi-user.target
